# Dotfiles/systemd/.config/systemd/user/home-assistant.service
[Unit]
Description=ShaderWebServer
After=network.target

[Service]
ExecStart=http-cmd

[Install]
WantedBy=multi-user.target

