# dotfiles/systemd/.config/systemd/user/home-assistant.service
[Unit]
Description=ShaderWebServer
After=network.target 

[Service]
Type=simple
Restart=always
ExecSearchPath=/home/canus/bin/:/home/server/bin/:/usr/bin/
ExecStart=http-cmd

[Install]
WantedBy=default.target

